# yopainstalld/Package/DEBIAN/postinst
#!/bin/sh

if [[ $1 == upgrade ]]; then
	/bin/launchctl unload /Library/LaunchDaemons/zorro.yopainstalld.plist
fi

chown root:wheel /usr/bin/__PACKAGENAME__
chmod 755 /usr/bin/__PACKAGENAME__
chown root:wheel /Library/LaunchDaemons/zorro.yopainstalld.plist
chmod 644 /Library/LaunchDaemons/zorro.yopainstalld.plist

/bin/launchctl load /Library/LaunchDaemons/zorro.yopainstalld.plist

exit 0
